# docs/reference/OpenViking/third_party/agfs/agfs-server/agfs-server.service
[Unit]
Description=AGFS Server - A General File System
Documentation=https://github.com/c4pt0r/agfs
After=network.target

[Service]
Type=simple
User=%USER%
Group=%GROUP%
ExecStart=%INSTALL_DIR%/agfs-server -c /etc/agfs.yaml
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
